# limnoria/usr/share/install/supybot@local.service.d/10-network.conf
[Service]
ExecStartPost=-/usr/bin/docker network create irc
ExecStartPost=/bin/sh -c "while ! /usr/bin/docker network \
connect irc supybot-local; do sleep 1; done"
