# toolbox/include/etc/profile.d/x-aliases.sh
alias ll='ls -alFh --group-directories-first'
alias la='ls -A'
alias l='ls -CF'
alias rm='rm -i'
alias cp='cp -i'

alias c='coredumpctl'
alias j='journalctl'
alias m='machinectl'
alias n='networkctl'
alias s='systemctl'
alias e='etcdctl'
alias k='kubectl'
alias h='chroot-host'
