# toolbox/etc/profile.d/x-aliases.sh
alias ll='ls -alFh --group-directories-first'
alias la='ls -A'
alias l='ls -CF'
alias rm='rm -i'
alias cp='cp -i'

alias c='chroot-host coredumpctl'
alias j='chroot-host journalctl'
alias m='chroot-host machinectl'
alias n='chroot-host networkctl'
alias s='chroot-host systemctl'

alias a='atomic'
alias e='etcdctl'
alias k='kubectl'
