# toolbox/usr/share/install/profile.d/docker-aliases.sh
../../../../etc/profile.d/docker-aliases.sh